# templates/etcd/k3s-etcd-defrag.service.j2
[Unit]
Description=Run Kubernetes ETCD defragment
Requires=k3s.service
After=k3s.service

[Service]
Type=oneshot
ExecStart={{ install.os.usr_bin }}//etcd-defrag.sh
